# ops/rehearsal/baby-quirt-host-certification.service
[Unit]
Description=Baby Quirt disposable host certification
After=local-fs.target
SuccessAction=poweroff
FailureAction=poweroff-force

[Service]
Type=oneshot
Environment=BABY_QUIRT_CERTIFICATION_PLAN=/run/baby-quirt-certification/input/plan.json
ExecStart=/opt/node-v24.18.0-linux-x64/bin/node /usr/local/libexec/baby-quirt-host-certification.mjs
StandardInput=null
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=96h
UMask=0077

[Install]
WantedBy=multi-user.target
